# scripts/kiosk/kiosk.service
[Unit]
Description=Kiosk Browser (ForenSys Dashboard)
# Arranca despues de la red y del gestor de seat (seatd)
After=network.target seatd.service forensys.service
Wants=network.target seatd.service
Requires=forensys.service

[Service]
User=ciber-admin
WorkingDirectory=/home/ciber-admin

# seatd crea el directorio /run/user/<uid> para el runtime de Wayland
RuntimeDirectory=kiosk
RuntimeDirectoryMode=0700
Environment="XDG_RUNTIME_DIR=/run/kiosk"

# Decirle a libseat que use seatd en lugar de logind
# (evita el error "Interactive authentication required")
Environment="LIBSEAT_BACKEND=seatd"

# Deshabilitar dispositivos de entrada virtuales que causan problemas
Environment="WLR_LIBINPUT_NO_DEVICES=1"

# NO usar StandardInput=tty ni TTYPath para no bloquear la consola
StandardOutput=journal
StandardError=journal

ExecStart=/usr/bin/cage -d -s -- \
  /usr/bin/chromium \
    --kiosk \
    --noerrdialogs \
    --disable-infobars \
    --incognito \
    --no-sandbox \
    --disable-dev-shm-usage \
    --ozone-platform=wayland \
    http://127.0.0.1:5000

# Si falla, esperar 10s antes de reintentar (evita loops rapidisimos)
Restart=on-failure
RestartSec=10s

[Install]
WantedBy=multi-user.target
