# scripts/kiosk/kiosk.service
[Unit]
Description=Kiosk Browser (ForenSys Dashboard)
After=network.target seatd.service forensys.service
Wants=network.target seatd.service
Requires=forensys.service

[Service]
User=ciber-admin
WorkingDirectory=/home/ciber-admin

RuntimeDirectory=kiosk
RuntimeDirectoryMode=0700
Environment="XDG_RUNTIME_DIR=/run/user/1000"

# Forzar card1 = RP1 que maneja HDMI en Raspberry Pi 5
Environment="WLR_DRM_DEVICES=/dev/dri/card1"

# Usar seatd en lugar de logind
Environment="LIBSEAT_BACKEND=seatd"

# Sin dispositivos de entrada virtuales
Environment="WLR_LIBINPUT_NO_DEVICES=1"

# Forzar renderizado por software (pixman) — evita pantalla azul por crash EGL/GPU en RPi5
Environment="WLR_RENDERER=pixman"

StandardOutput=journal
StandardError=journal

# Esperar a que la web app esté lista antes de lanzar el navegador
ExecStartPre=/bin/sleep 5

ExecStart=/usr/bin/cage -s -- \
  chromium-browser \
    --kiosk \
    --noerrdialogs \
    --disable-infobars \
    --incognito \
    --no-sandbox \
    --disable-dev-shm-usage \
    --ozone-platform=wayland \
    --disable-gpu \
    --disable-software-rasterizer \
    --ignore-gpu-blocklist \
    --use-gl=egl \
    http://127.0.0.1:5000

Restart=on-failure
RestartSec=10s

[Install]
WantedBy=multi-user.target
